# deploy/tickstone-sync.service
[Unit]
Description=TickStone Bluetooth history sync
Requires=bluetooth.service tickstone-ble-connection.service
After=bluetooth.service tickstone-ble-connection.service

[Service]
Type=simple
User=allangamal
WorkingDirectory=/home/allangamal/TickStone
Environment=HOME=/home/allangamal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/allangamal/TickStone/.venv/bin/python /home/allangamal/TickStone/tools/tickstone_ble_sync.py --watch --output /home/allangamal/.local/share/tickstone/logs.jsonl --database /home/allangamal/.local/share/tickstone/tickstone.sqlite3
Restart=always
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=read-only
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_BLUETOOTH
ReadWritePaths=/home/allangamal/.local/share/tickstone

[Install]
WantedBy=multi-user.target
